# ops/systemd/codexnext-control.service.example
[Unit]
Description=CodexNext Relay Control Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__CODEXNEXT_USER__
WorkingDirectory=__CODEXNEXT_ROOT__
EnvironmentFile=/etc/codexnext/control.env
ExecStart=__CODEXNEXT_ROOT__/scripts/ops/run-control.sh
Restart=always
RestartSec=3
TimeoutStopSec=20

[Install]
WantedBy=multi-user.target
